# package/rfm12_server/files/rfm12.init
#!/bin/sh /etc/rc.common
START=100

start()
{
	/sbin/rfm12_server > /dev/null &
}

stop()
{
	killall rfm12_server
}
